# .config/autostart/screensaver-settings.desktop
[Desktop Entry]
Exec=xset s 3000 dpms 0 0 0
Type=Application
Name=Screensaver
Comment=Set screensaver timeouts
