# .config/autostart/screensaver-settings.desktop
[Desktop Entry]
Exec=xset s 360 dpms 480 0 0
Type=Application
Name=Screensaver
Comment=Set screensaver timeouts
